# crontab
0       12      *       *       *       /usr/bin/python /home/poppetode/maintenance.py > /dev/null
@reboot /usr/bin/python /home/poppetode/maintenance.py noreboot > /dev/null

